# deploy/systemd/subscription-manager-backup.timer
[Unit]
Description=Daily Subscription Manager backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=15m
Unit=subscription-manager-backup.service

[Install]
WantedBy=timers.target
